# .cron
# begin dotfile jobs
# DO NOT REMOVE THE ABOVE LINE
*/15 * * * * (crontab -l | sed -e '/^# begin dotfile jobs$/,/^# end dotfile jobs$/d'; cat $HOME/.cron) | crontab -
*/5 * * * * "$HOME/bin/tmux-prune-local-sessions.sh"
# DO NOT REMOVE THE BELOW LINE
# end dotfile jobs
